# runtime/etc/init/berry.conf
description "berry agent to retrieve application credentials"

start on started network-services
stop on shutdown

respawn

exec berry /meta/credentials
