# runtime/etc/init/configure-ssh.conf
description "Configure OpenSSH server"
start on (cloud-config and (starting ssh))

task

exec /opt/zalando/bin/configure-ssh
