# runtime/etc/init/ecr_login.conf
# after docker is initialized
start on started taupage

# this is only a onetime task
task

# the actual script
script
    `aws ecr get-login --region us-east-1` >/dev/console 2>&1
end script
