# runtime/etc/init/zalando.conf
# after docker is initialized
start on started docker

# this is only a onetime task
task

# the actual script
script
    /opt/taupage/init.sh 2>&1 | logger -t taupage-init
end script
