# runtime/etc/cron.hourly/logrotate_auditlog
#!/bin/bash

logrotate -f /etc/logrotate.d/auditlog
