# build/files
/etc/apt/sources.list.d/zalando-ami.list root:root 0444

/etc/audit/audit.rules root:root 0400
/etc/audisp/plugins.d/syslog.conf root:root 0400

/etc/cloud/cloud.cfg root:root 0644

/etc/cron.hourly/logrotate_auditlog root:root 0755

/etc/logrotate.d/auditlog root:root 0644

/etc/sudoers root:root 0440

/etc/newrelic newrelic:newrelic 0750
/etc/newrelic/nrsysmond.cfg newrelic:newrelic 0640

/etc/rsyslog.d/keys root:root 0755
/etc/rsyslog.d/keys/ca.d root:root 0755

/etc/ssh-access-granting-service.yaml root:root 0444

/etc/taupage.yaml root:root 0444

/opt/taupage/init.d/10-prepare-disks.py root:root 0700
/opt/taupage/init.sh root:root 0700
/opt/taupage/runtime/Docker root:root 0700

/usr/lib/python2.7/dist-packages/cloudinit/handlers/zalando_ami_config.py root:root 0444

/tmp/delete-ubuntu-user.sh root:root 0700
/tmp/delete-ubuntu-user-wrapper.sh root:root 7755
