# runtime/etc/cron.d/log_agent_mem
*/1 * * * * root /opt/taupage/bin/log_agent_mem.sh
