# runtime/etc/cron.d/mon-put-instance-data-cloudwatch.deactivated
*/5 * * * * root /opt/taupage/cloudwatch/aws-scripts-mon/mon-put-instance-data.pl --mem-util --disk-space-util --from-cron
